# src/backend/supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=flask --app manage:app run --debug --host=0.0.0.0 --port=8080
directory=/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery]
command=celery -A app.factory:celery_app worker --loglevel=info
directory=/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
